# ansible/roles/node-opensearch/files/lib/systemd/system/caddy.service.d/10-opensearch.conf
[Service]
User=opensearch
Group=opensearch
